# auto-restart-services/systemd timer/ocsanitizeservice.service
[Unit]
Description=Opencast Service Sanitizer

[Service]
ExecStart=/root/scripts/sanitize_services.sh
